# modep-mod-midi-merger/modep-mod-midi-merger-0.1.0/debian/modep-mod-midi-merger.modep-mod-midi-merger.service
[Unit]
Description=MOD MIDI Merger
After=jack.service
Requires=jack.service
BindsTo=jack.service

[Service]
RemainAfterExit=yes
Environment=JACK_PROMISCUOUS_SERVER=jack
ExecStart=/usr/bin/jack_load mod-midi-merger
ExecStop=/usr/bin/jack_unload mod-midi-merger
User=modep
Group=modep
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
