# modep-common/modep-common-1.1.0/debian/modep-common.modep-update.path
[Unit]
Description=Modep Update Monitor path.

[Path]
PathExists=/tmp/modep-update
Unit=modep-update.service

[Install]
WantedBy=multi-user.target
